# crontab
0 7 * * 1-5 /app/go-email-fleury >> /var/log/cron.log 2>&1